# daemons/desktop/syrem.service.in
[Unit]
Description=Syrem Daemon
After=network-online.target

[Service]
Type=simple
ExecStart=%{INSTALL_BINS}/syremd --systemd-log
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
